# backend/docker/etc/s6-overlay/s6-rc.d/init-db-migration/run
#!/command/with-contenv bash
# shellcheck shell=bash

[ ${MIXMATCH_ENABLE_MIGRATE} == "true" ] && mixmatch admin migrate
